# bot/bot-tool.install/etc/systemd/system/bot-telegram.service
[Unit]
Description=Bot-Telegram service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=anselmo
ExecStart=/usr/bin/bot.py
